# mcp/jamong-mcp.service
[Unit]
Description=Jamong Skills MCP Server
After=network.target

[Service]
Type=simple
WorkingDirectory=__REPO_DIR__
ExecStart=__VENV_DIR__/bin/python __REPO_DIR__/mcp/server.py
Restart=always
RestartSec=5
Environment=PORT=8000

[Install]
WantedBy=multi-user.target
